# .config/systemd/user/battery-monitor.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec %h/bin/battery-monitor"
