# .config/systemd/user/battery-monitor.service
[Service]
ExecStart=/bin/sh -c "exec %h/bin/battery-monitor"
TimeoutStartSec=5s
Restart=always
RestartSec=60
