# .config/systemd/user/battery-notify.service
[Service]
ExecStart=/bin/sh -c "exec %h/.cargo/bin/battery-notify"
