# .config/systemd/user/clipmenud.service
[Service]
ExecStart=/usr/bin/clipmenud
Restart=always
