# .config/systemd/user/dwm.service
[Service]
Type=oneshot
ExecStart=/usr/bin/dwm
