# .config/systemd/user/fcitx.service
[Service]
ExecStart=/usr/bin/fcitx -D
Restart=on-failure
