# .config/systemd/user/firefox.service
[Service]
ExecStart=/usr/bin/firefox
Restart=always
Environment="NSS_DEFAULT_DB_TYPE=sql"
