# .config/systemd/user/firefox.service
[Service]
ExecStart=firefox
Restart=always
Environment="NSS_DEFAULT_DB_TYPE=sql"
