# .config/systemd/user/fontdirs.service
[Service]
Type=oneshot
ExecStart=xset +fp /usr/share/fonts/local
ExecStart=xset fp rehash
