# .config/systemd/user/fontdirs.service
[Service]
Type=oneshot
ExecStart=/usr/bin/xset +fp /usr/share/fonts/local
ExecStart=/usr/bin/xset fp rehash
