# .config/systemd/user/habitual.service
[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=pyenv exec habitual
Restart=always
