# .config/systemd/user/inbox-reminder.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec %h/bin/inbox-reminder"
