# .config/systemd/user/inbox-reminder.timer
[Timer]
OnUnitActiveSec=120min
OnActiveSec=120min
