# .config/systemd/user/irssi.service
[Service]
ExecStart=/bin/sh -c "exec %h/bin/start-irssis"
Restart=on-failure
