# .config/systemd/user/irssi@.service
[Service]
ExecStart=irssi --home %I
Restart=always
