# .config/systemd/user/irssi.service
[Unit]
ConditionPathExists=/usr/bin/irssi

[Service]
ExecStart=/bin/sh -c "exec %h/bin/start-irssis"
Restart=on-failure
