# .config/systemd/user/jupyter-notebook.service
[Service]
WorkingDirectory=%h/Notebooks
ExecStart=/bin/sh -c "exec %h/bin/jpn"
Restart=always
