# .config/systemd/user/mbsync.service
[Unit]
After=network.target

[Service]
ExecStart=/bin/sh -c "exec %h/bin/run-mbsyncs"
TimeoutStartSec=1min30s
Restart=always
RestartSec=5
