# .config/systemd/user/mpdscribble.service
[Unit]
Description=MPD last.fm scrobbler

[Service]
ExecStart=/usr/bin/mpdscribble --no-daemon

[Install]
WantedBy=default.target
