# .config/systemd/user/mute-all-mics.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec %h/bin/mute-all-mics"
