# .config/systemd/user/mutt.service
[Unit]
ConditionPathExists=/usr/bin/mutt

[Service]
ExecStart=/bin/sh -c "exec %h/bin/start-mutts"
Restart=on-failure
