# .config/systemd/user/no-dpms.service
[Service]
Type=oneshot
ExecStart=/usr/bin/xset dpms 0 0 0
