# .config/systemd/user/no-screensaver.service
[Service]
Type=oneshot
ExecStart=/usr/bin/xset s 0
