# .config/systemd/user/offlineimap.service
[Unit]
Description=Synchronise local and remote mails

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap
