# .config/systemd/user/offlineimap.service
[Unit]
Description=OfflineIMAP
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap -o -q -u quiet
TimeoutStartSec=1min30s
