# .config/systemd/user/psi-notify.service
[Service]
ExecStart=/usr/bin/psi-notify
Restart=always
