# .config/systemd/user/redshift-auto.service
[Service]
Type=exec
ExecStart=/bin/sh -c "exec %h/bin/redshift-auto"
