# .config/systemd/user/spotify.service
[Service]
ExecStart=/bin/sh -c "exec %h/bin/spotify"
Restart=on-failure
