# .config/systemd/user/tmux.service
[Service]
Type=oneshot
ExecStart=/usr/bin/tmux start-server
