# .config/systemd/user/todo.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec %h/bin/maybe-todo-reminder"
