# .config/systemd/user/todo.service
[Service]
Type=oneshot
ExecStart=/usr/bin/nota-todo-reminder
