# .config/systemd/user/todo.service
[Service]
Type=oneshot
ExecStart=nota-todo-reminder
