# .config/systemd/user/todo.timer
[Timer]
OnUnitActiveSec=30min
OnActiveSec=30min
