# .config/systemd/user/upgrade-reminder.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "exec %h/bin/upgrade-reminder"
