# .config/systemd/user/urxvtd.service
[Unit]
Description=rxvt-unicode terminal daemon

[Service]
ExecStart=/usr/bin/urxvtd -q -o

[Install]
WantedBy=default.target
