# .config/systemd/user/urxvtd.service
[Service]
ExecStart=/usr/bin/urxvtd -q -o
