# .config/systemd/user/vboxclient.service
[Unit]
ConditionVirtualization=oracle

[Service]
Type=oneshot
ExecStart=/usr/bin/VBoxClient-all
Restart=on-failure
