# .config/systemd/user/xautolock.service
[Service]
ExecStart=/usr/sbin/xautolock -time 1 -locker "%h/bin/lock-maybe /bin/slock"
Restart=always
