# .config/systemd/user/xautolock.service
[Service]
ExecStart=/usr/sbin/xautolock -time 1 -locker /bin/slock
Restart=always
