# .config/systemd/user/xbanish.service
[Service]
ExecStart=/usr/bin/xbanish -i Mod4
Restart=on-failure
