# .config/systemd/user/xinput-disable@.service
[Service]
Type=oneshot
ExecStart=/usr/bin/xinput-toggle -r %I -d
