# .config/systemd/user/xob.service
[Service]
ExecStart=/bin/sh -c "exec %h/bin/run-xob"
Restart=always
